# supervisord.conf
; ═══════════════════════════════════════════════════════════════
; supervisord.conf — Monica Medical AI Server
;
; 用于 Docker 容器内（由 CMD 调用：supervisord -c supervisord.conf -n）
; 裸机部署请使用 bare-deploy.sh，它会自动生成专属的 supervisor 配置
; 写入 /etc/supervisor/conf.d/monica.conf
; ═══════════════════════════════════════════════════════════════

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/dev/null          ; 容器内日志走 stdout，不写文件
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
silent=false
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── FastAPI / Uvicorn ──────────────────────────────────────────
; Docker 模式：监听 0.0.0.0（由 nginx 容器做反向代理）
; 裸机模式：bare-deploy.sh 生成的配置监听 127.0.0.1（由宿主 nginx 代理）
[program:fastapi]
command=uvicorn app.main:app
    --host 0.0.0.0
    --port 8000
    --workers 1
    --loop uvloop
    --log-level info
    --no-access-log
directory=%(ENV_APP_DIR)s
user=%(ENV_APP_USER)s
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="%(ENV_APP_DIR)s"

; ── ARQ Worker ─────────────────────────────────────────────────
[program:arq_worker]
command=arq app.workers.arq_worker.WorkerSettings
directory=%(ENV_APP_DIR)s
user=%(ENV_APP_USER)s
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="%(ENV_APP_DIR)s"

[group:monica]
programs=fastapi,arq_worker
