# supervisord.conf
; ═══════════════════════════════════════════════════════════════
; supervisord.conf — Monica Medical AI Server
;
; Docker 容器内由 CMD 调用：supervisord -c supervisord.conf -n
; 裸机部署时修改 directory / environment 中的路径即可
; ═══════════════════════════════════════════════════════════════

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/dev/null          ; 容器内日志走 stdout，不写文件
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false             ; Docker CMD 模式下 supervisord -n 已处理前台
silent=false
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── FastAPI / Uvicorn ──────────────────────────────────────────
[program:fastapi]
command=uvicorn app.main:app
    --host 0.0.0.0
    --port 8000
    --workers 1
    --loop uvloop
    --log-level info
    --no-access-log
directory=%(ENV_APP_DIR)s
user=%(ENV_APP_USER)s
autostart=true
autorestart=true
startretries=5
startsecs=5
; 日志直接写到 stdout/stderr，由 Docker/systemd 收集
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="%(ENV_APP_DIR)s"

; ── ARQ Worker ────────────────────────────────────────────────
[program:arq_worker]
command=arq app.workers.arq_worker.WorkerSettings
directory=%(ENV_APP_DIR)s
user=%(ENV_APP_USER)s
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
redirect_stderr=false
environment=PYTHONPATH="%(ENV_APP_DIR)s"

[group:monica]
programs=fastapi,arq_worker
